# ansible/templates/microshift-lvm-loop.service.j2
# Ansible-managed — attach the sparse LVM image as a loop device before LVM starts.
[Unit]
Description=Attach MicroShift LVM sparse image ({{ microshift_loop_dev }})
Documentation=file:///home/swares/lab/homelab/homelab/docs/ARCHITECTURE.md
After=mnt-nvme0n1p2.mount
Requires=mnt-nvme0n1p2.mount
Before=lvm2-monitor.service microshift.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c 'losetup {{ microshift_loop_dev }} {{ microshift_lvm_file }} || losetup {{ microshift_loop_dev }} | grep -q microshift-lvm.img'
ExecStart=/sbin/vgchange -ay {{ microshift_vg }}
ExecStop=/sbin/vgchange -an {{ microshift_vg }}
ExecStop=/usr/sbin/losetup -d {{ microshift_loop_dev }}

[Install]
WantedBy=multi-user.target
