# ansible/ansible/templates/backup-etcd.service.j2
[Unit]
Description=k3s state backup to cold tier (8TB)
RequiresMountsFor={{ cold_primary_mount }}

[Service]
Type=oneshot
# k3s single-node uses SQLite (state.db) — no embedded etcd.
# sqlite3 .backup does an online hot-copy (safe while k3s is running).
# Keeps the last 7 copies; prunes older ones.
ExecStart=/bin/bash -c '\
  DB=/var/lib/rancher/k3s/server/db/state.db && \
  DEST={{ cold_primary_mount }}/k3s-etcd-snapshots && \
  SNAP="$DEST/state-$(date +%%F_%%H%%M).db" && \
  mkdir -p "$DEST" && \
  sqlite3 "$DB" ".backup $SNAP" && \
  ls -1t "$DEST"/state-*.db | tail -n +8 | xargs -r rm -f'
